# usr/lib/systemd/system/docker-cleanup.timer
[Unit]
Description=ostree-docker-cleanup timer

[Timer]
OnBootSec=15m
OnUnitActiveSec=8h
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
