# usr/lib/systemd/system/var-opt-persist.mount
[Mount]
What=/dev/disk/by-id/scsi-0QEMU_QEMU_HARDDISK_drive-scsi2-part1
Where=/var/opt/persist
Type=ext4

[Install]
WantedBy=local-fs.target
