# usr/lib/systemd/system/zfs-snap@.service
[Unit]
Description=zfs-snap %i

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/zfs-snap.sh %i
