# usr/lib/systemd/system/zfs-snap@.timer
[Unit]
Description=zfs-snap timer

[Timer]
RandomizedDelaySec=5m
OnBootSec=15min
OnUnitActiveSec=4h

[Install]
WantedBy=timers.target