# usr/lib/systemd/system/zpool-scrub@.timer
[Unit]
Description=zpool scrub timer

[Timer]
OnCalendar=Mon *-*-01..07 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
