# usr/etc/systemd/system/target.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/targetctl restore /etc/target/masterconfig.json
